# instaclaw/scripts/browser-relay-server/browser-relay-server.service
[Unit]
Description=InstaClaw Browser Relay Server
After=network.target openclaw-gateway.service
Wants=openclaw-gateway.service
# Crash-loop protection — StartLimit directives belong in [Unit] (systemd >=240)
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
ExecStart=/usr/bin/node /home/openclaw/scripts/browser-relay-server.js
Environment=HOME=/home/openclaw
Environment=PATH=/home/openclaw/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/usr/bin:/bin
Environment=RELAY_PORT=18792
Environment=RELAY_BIND=127.0.0.1
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
